# rootfs/etc/s6-overlay/s6-rc.d/otbr-agent/run
#!/usr/bin/with-contenv bash
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# OpenThread BorderRouter start script
# ==============================================================================

# shellcheck source=./openthread_border_router/rootfs/etc/s6-overlay/scripts/otbr-agent-common
. /etc/s6-overlay/scripts/otbr-agent-common

declare trel_if
declare backbone_if
declare device
declare baudrate
declare flow_control
declare migrate_flow_control
declare otbr_log_level
declare otbr_log_level_int
declare otbr_rest_listen
declare otbr_rest_listen_port

backbone_if=${BACKBONE_IF:-eth0}
trel_if=${TREL_IF}
baudrate=${BAUDRATE}
device=${DEVICE}

# Backbone network interface can be configured by user or auto-detected by looking for the primary network interface in the system.
# If no primary network interface is found, the app will exit with an error and ask the user to configure a backbone network interface in the app configuration.
if [ -z "${backbone_if}" ] ; then
    # If no backbone interface is found, exit with error
    echo "ERROR: No primary network interface found. Please configure a backbone network interface with the BACKBONE_IF environment variable"
    exit 1
fi

if [ -z ${trel_if} ]; then
    echo "WARNING: No primary trel interface found! Using ${backbone_if} interface"
    trel_if=${backbone_if}
fi

# Check if IPv6 forwarding is enabled
if [ "$(cat /proc/sys/net/ipv6/conf/all/forwarding)" != "1" ]; then
    echo "WARNING: IPv6 routing/forwarding is not enabled! Make sure the Home Assistant host has IPv6 forwarding enabled."
fi

flow_control=""
if [ x"${FLOW_CONTROL}" = x"1" ] ; then
    flow_control="&uart-flow-control"
    migrate_flow_control="hardware"
else
    flow_control="&uart-init-deassert"
    migrate_flow_control="none"
fi

otbr_log_level=${OTBR_LOG_LEVEL:-info}
case "${otbr_log_level}" in
    debug)
        otbr_log_level_int="7"
        ;;
    info)
        otbr_log_level_int="6"
        ;;
    notice)
        otbr_log_level_int="5"
        ;;
    warning)
        otbr_log_level_int="4"
        ;;
    error)
        otbr_log_level_int="3"
        ;;
    critical)
        otbr_log_level_int="2"
        ;;
    alert)
        otbr_log_level_int="1"
        ;;
    emergency)
        otbr_log_level_int="0"
        ;;
    *)
        echo "NOTICE: Unknown OTBR_LOG_LEVEL: ${otbr_log_level}. Using level info"
        otbr_log_level_int=6
        ;;
esac

# shellcheck disable=SC2015
mkdir -p /data/thread && ln -sft /var/lib /data/thread
rc=$? ; if [ x"$rc" != x"0" ] ; then
    echo "ERROR: Could not create directory /var/lib/thread to store Thread data."
    exit $rc
fi

# We compile the OTBR with firewall support, so otbr-agent tries to update the
# ipsets. Therefor, create ipsets always to avoid errors from otbr-agent. Just
# the ipsets  won't have an effect in practice when the firewall is disabled.
ipset create -exist otbr-ingress-deny-src hash:net family inet6
ipset create -exist otbr-ingress-deny-src-swap hash:net family inet6
ipset create -exist otbr-ingress-allow-dst hash:net family inet6
ipset create -exist otbr-ingress-allow-dst-swap hash:net family inet6

ip6tables -N "${otbr_forward_ingress_chain}"
ip6tables -I FORWARD 1 -o "${thread_if}" -j "${otbr_forward_ingress_chain}"

ip6tables -N "${otbr_forward_egress_chain}"
ip6tables -I FORWARD 2 -i "${thread_if}" -j "${otbr_forward_egress_chain}"

if [ x"${FIREWALL}" = x"1" ] ; then
    echo "INFO: Setup OTBR firewall..."

    ip6tables -A "${otbr_forward_ingress_chain}" -m pkttype --pkt-type unicast -i "${thread_if}" -j DROP
    ip6tables -A "${otbr_forward_ingress_chain}" -m set --match-set otbr-ingress-deny-src src -j DROP
    ip6tables -A "${otbr_forward_ingress_chain}" -m set --match-set otbr-ingress-allow-dst dst -j ACCEPT
    ip6tables -A "${otbr_forward_ingress_chain}" -m pkttype --pkt-type unicast -j DROP
    ip6tables -A "${otbr_forward_ingress_chain}" -j ACCEPT

    ip6tables -A "${otbr_forward_egress_chain}" -j ACCEPT
else
    ip6tables -A "${otbr_forward_ingress_chain}" -j ACCEPT

    ip6tables -A "${otbr_forward_egress_chain}" -j ACCEPT
fi

if [ x"${NAT64}" = x"1" ] ; then
    # Mark Thread traffic in mangle
    iptables -t mangle -A PREROUTING -i "${thread_if}" -j MARK --set-mark "${otbr_fw_mark}"

    # MASQUERADE marked traffic
    iptables -t nat -A POSTROUTING -m mark --mark "${otbr_fw_mark}" -j MASQUERADE

    # NAT64 forward chain — jump unconditionally, filter inside
    iptables -N "${otbr_forward_nat64_chain}"
    iptables -I FORWARD 1 -j "${otbr_forward_nat64_chain}"

    # Forward marked traffic
    iptables -A "${otbr_forward_nat64_chain}" -m mark --mark "${otbr_fw_mark}" -o "${backbone_if}" -j ACCEPT
    iptables -A "${otbr_forward_nat64_chain}" -m conntrack --ctstate ESTABLISHED,RELATED -i "${backbone_if}" -o "${thread_if}" -j ACCEPT
fi

otbr_rest_listen="::"
otbr_rest_listen_port="${OTBR_REST_PORT:-8081}"

# If user port is not set, listen on local interface only
if [ -z "${otbr_rest_listen_port}" ] ; then
    otbr_rest_listen_port="8081"
elif [ "${otbr_rest_listen_port}" != "8081" ]; then
    echo "WARNING: Custom OpenThread REST API port is not supported. Using 8081."
    otbr_rest_listen_port="8081"
fi

# Store REST API listen information for check script
echo "${otbr_rest_listen}" > /tmp/otbr-agent-rest-api
echo "${otbr_rest_listen_port}" >> /tmp/otbr-agent-rest-api

# Migrate OTBR settings to new adapter if needed
echo "INFO: Migrating OTBR settings if needed..."
python3 /usr/local/bin/migrate_otbr_settings.py \
    --adapter "${device}" \
    --baudrate "${baudrate}" \
    --flow-control "${migrate_flow_control}" \
    --data-dir /data/thread/

echo "INFO: Starting otbr-agent..."
# shellcheck disable=SC2086
exec s6-notifyoncheck -d -s 300 -w 300 -n 0 stdbuf -oL \
    "/usr/sbin/otbr-agent" -I ${thread_if} -B "${backbone_if}" \
        --rest-listen-address "${otbr_rest_listen}" \
        --rest-listen-port "${otbr_rest_listen_port}" \
        -d${otbr_log_level_int} -v -s \
        "spinel+hdlc+uart://${device}?uart-baudrate=${baudrate}${flow_control}" \
        "trel://${trel_if}"
